# app/docker/supervisord.conf
[supervisord]
nodaemon=true
[program:app]
command=php /www/run.php